# teletraan-deployment/teletraan-ha.service
# ============================================================
# Systemd unit for the HA failover watcher
# Copy to: /etc/systemd/system/teletraan-ha.service
# Then:  systemctl daemon-reload && systemctl enable --now teletraan-ha
# ============================================================

[Unit]
Description=Teletraan I HA Failover Watcher
After=network-online.target teletraan-agent.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/bash /opt/teletraan/ha-failover.sh 192.168.1.20 7070
# ↑ Change 192.168.1.20 to your TrueNAS IP
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teletraan-ha

[Install]
WantedBy=multi-user.target
